# deploy/chronosync.service
[Unit]
Description=ChronoSync FastAPI
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/srv/chronosync
EnvironmentFile=/srv/chronosync/.env
ExecStart=/srv/chronosync/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
